# deploy/porto-web.service
# systemd unit for the Next.js frontend (standalone server).
# Install to: /etc/systemd/system/porto-web.service
# Then: sudo systemctl daemon-reload && sudo systemctl enable --now porto-web

[Unit]
Description=Porto Next.js Web
After=network.target porto-api.service
Wants=porto-api.service

[Service]
Type=simple
User=porto
Group=porto
WorkingDirectory=/var/www/porto/nextjs
Environment=NODE_ENV=production
Environment=PORT=3000
# Uses the standalone build output (next.config.js has output: "standalone").
ExecStart=/usr/bin/node /var/www/porto/nextjs/.next/standalone/server.js
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
